# AC100/Brightness/etc/udev/rules.d/99-backlight.rules
KERNEL=="backlight" SUBSYSTEM=="backlight" RUN+="/sbin/setup_backlight.sh"
